# dot_config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland Session

PartOf=graphical.target
After=graphical.target

Requires=graphical-session.target
After=graphical.target
