# dot_config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland Session
Requires=graphical-session.target
