# dot_config/systemd/user/waybar-restart.path
[Unit]
Description=Watch for Waybar config changes

[Path]
PathChanged=%h/.config/waybar/
PathChanged=%h/.config/waybar/config.jsonc
PathChanged=%h/.config/waybar/style.css

[Install]
WantedBy=graphical-session.target
