# dot_config/systemd/user/waybar-restart.service
[Unit]
Description=Restart Waybar on config change

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart waybar.service
